# forensic-cpa-ai.service
[Unit]
Description=Forensic CPA AI - Financial Private Investigator
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/user/Forensic_CPA_AI
Environment="PORT=5000"
Environment="FLASK_ENV=production"
ExecStart=/usr/bin/python3 /home/user/Forensic_CPA_AI/app.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
